# config/systemd/user/newsboat.timer
[Unit]
Description=Newsboat reload timer

[Timer]
OnBootSec=1m
OnUnitActiveSec=10m
Unit=newsboat.service

[Install]
WantedBy=timers.target
